# build/etcd/host/lib/systemd/system/etcd-proxy.service
[Unit]
Description=etcd proxy

[Service]
ExecStart=/usr/local/bin/etcd -proxy on

Type=notify
Restart=always

[Install]
WantedBy=multi-user.target
