# build/etcd/host/lib/systemd/system/etcd-proxy.service
[Unit]
Description=etcd proxy

[Service]
ExecStart=/usr/local/bin/etcd \
  -proxy on \
  -initial-cluster-token etcd-cluster-01 \
  -advertise-client-urls http://192.168.3.250:2379 \
  -listen-client-urls  http://192.168.3.117:2379,http://127.0.0.1:2379 \
  -initial-cluster etcd-01=http://192.168.3.250:2380,etcd-02=http://192.168.3.251:2380,etcd-03=http://192.168.3.252:2380 \

Type=notify
Restart=always

[Install]
WantedBy=multi-user.target
